# templates/jekyll_build.cron
0  * * * * root /usr/local/sbin/jekyll-build-{{ jekyll_site_name }} >/dev/null
30 * * * * root /usr/local/sbin/jekyll-build-{{ jekyll_site_name }} >/dev/null
