# templates/jekyll_build.cron
0  * * * * /usr/local/sbin/jekyll-build-{{ jekyll_site_name }}
30 * * * * /usr/local/sbin/jekyll-build-{{ jekyll_site_name }}
